# frankenstein-web.service
[Unit]
Description=Frankenstein Web Administration Dashboard
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/frankenstein-web.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
